# garden_journal/rootfs/etc/services.d/garden_journal/run
#!/usr/bin/with-contenv bashio
bashio::log.info "Starting Garden Journal..."
cd /app
exec node server.js
